# progs/glassfish_start_database
### BEGIN INIT INFO
# Provides:          glassfish_startup
# Short-Description: Stop tomcat and start glassfish
# Description:       This script stop the tomcat server and start java database and administration domain
### END INIT INFO

#!/bin/bash

echo 'Démarrage de la base de donnée Java...'
sudo /opt/glassfish3/glassfish/bin/asadmin start-database
