# sql/glassfish_startup
### BEGIN INIT INFO
# Provides:          glassfish_startup
# Short-Description: Stop tomcat and start glassfish
# Description:       This script stop the tomcat server and start java database and administration domain
### END INIT INFO

#!/bin/bash

cd /opt/glassfish3/glassfish/bin/
echo 'Stopping tomcat instance...'
sudo /etc/init.d/tomcat7 stop
echo 'Stopping database...'
#jdbc=$(sudo /opt/glassfish3/glassfish/bin/asadmin stop-database)
#echo $jdbc
echo 'Stopping administration domain...'
dom=$(sudo /opt/glassfish3/glassfish/bin/asadmin stop-domain administration)
echo $dom
echo 'Restarting administration domain...'
dom=$(sudo /opt/glassfish3/glassfish/bin/asadmin start-domain administration)
echo $dom
echo 'Program will continue in 60 secondes...'
#sleep 1m
echo 'Restarting database...'
#jdbc=$(sudo /opt/glassfish3/glassfish/bin/asadmin start-database)
#echo $jdbc
